# debs/com.amirramos.theos/etc/profile.d/theos.sh
export THEOS=/var/theos
export PATH=$THEOS/bin:$PATH
